# crontab
@daily flock -n /tmp/cron-flexget.flock /usr/bin/flexget --cron execute
